# services/cron/crontab
# Every 1 minute: Process Message Queue
* * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/process-queue" > /dev/null

# Every 10 minutes: Rescue Abandoned Conversations
*/10 * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/rescue" > /dev/null

# Every 5 minutes: Remind Voice Source
*/5 * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/remind-voice" > /dev/null

# Every 1 minute: Remind Lead Provider (Urgent)
* * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/remind-leads" > /dev/null

# Every 2 minutes: Recover Failed Webhooks
*/2 * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/recovery" > /dev/null

# Every 1 minute: Process RVC Voice Conversion Jobs
* * * * * curl -s "https://main.die6qcz48fek8.amplifyapp.com/api/cron/process-voice" > /dev/null
